# .config/systemd/user/atuin-daemon.socket
[Unit]
Description=Unix socket activation for atuin daemon

[Socket]
ListenStream=%t/atuin.sock

[Install]
WantedBy=sockets.target

